# config/room-scanner.service
[Unit]
Description=3D Room Scanner Service
After=network.target multi-user.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%PROJECT_DIR%/ros2_ws
Environment="ROS_DOMAIN_ID=0"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/%USER%/.Xauthority"
ExecStartPre=/bin/sleep 10
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source %PROJECT_DIR%/ros2_ws/install/setup.bash && ros2 launch room_scanner scanner.launch.py'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
